# src/overlay/root/etc/systemd/system/kill-splash.service
[Unit]
Description=Kill Splash Screen

[Service]
ExecStart=/usr/bin/killall fbi
Type=oneshot

[Install]
WantedBy=getty.target
